# recipes-core/busybox/files/watchdog.service.in
[Unit]
Description=Watchdog daemon

[Service]
Type=forking
ExecStart=@base_sbindir@/watchdog -t 5 /dev/watchdog

[Install]
WantedBy=multi-user.target
